# supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:bamboo-agent]
command=java -jar /"%(ENV_bamboo_jar)s" "%(ENV_bamboo_url)s"/bamboo/agentServer
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
